# systemd/wellness.timer
[Unit]
Description=Run wellness.service every weekday at 9:50am Mountain Time

[Timer]
OnCalendar=Mon-Fri *-*-* 15:50:00
